# deploy/labber-agent.service
[Unit]
Description=LabberAirport Agent (sing-box control plane)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=labber
Group=labber
WorkingDirectory=/var/lib/labber-agent
EnvironmentFile=/etc/labber-agent/agent.env
ExecStart=/usr/local/bin/labber-agent -listen ${LABBER_LISTEN} -token ${LABBER_TOKEN} -data-dir ${LABBER_DATA_DIR}
Restart=on-failure
RestartSec=3
LimitNOFILE=1048576
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/labber-agent
PrivateTmp=true

[Install]
WantedBy=multi-user.target
